# docker/ruby265/servers/redmine/supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

#[program:init]
#command=/usr/local/bin/init.container
#autorestart=false
#startretries=0
#exitcodes=0



