# lamp/docker/rabbitmq/supervisor.conf
[supervisord]
#user=root
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:rabbitmq]
redirect_stderr=true
command=/usr/sbin/rabbitmq-server
process_name = rabbitmq
