# docker/primavera/18.8/primavera.conf
[program:init]
command=/usr/local/bin/init.container
autorestart=false
startretries=0
exitcodes=0
