# docker/db/maria/supervisor.conf
[program:sshd]
command=/usr/sbin/sshd -D

[program:mysql]
command=/usr/bin/mysqld_safe

[program:init]
command=/usr/local/bin/init.container
autorestart=false
startretries=0
exitcodes=0


