# docker/java/gitlab/server/supervisor.gitlab.conf
[program:gitlab]
command=/bin/bash -c "/usr/local/bin/init.install"
autorestart=false
startretries=0
exitcodes=0
startsecs=30
