# docker/php53/sites/kisnew/supervisor.conf
[program:cron]
#command=/bin/bash -c "exec /usr/sbin/crond -n"
command=/usr/sbin/crond
autostart=true
autorestart=true
startretries=50
