# docker/db/postgre/1c/supervisor.conf
[program:pgsql]
command=/usr/local/bin/postgres

[program:sshd]
command=/usr/sbin/sshd -D
