# docker/wowza/supervisor.conf
[program:sshd]
command=/usr/sbin/sshd -D

[program:wowza]
command=/usr/local/bin/runWowza
