# services/docker/mail/dovecot/supervisor.conf
[program:sshd]
command=/usr/sbin/sshd -D

[program:postfix]
command=/usr/sbin/postfix start

[program:dovecot]
command=/usr/sbin/dovecot
