# docker/turn/supervisor.conf
[supervisord]
#user=root
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:turn]
##redirect_stderr=true
command=/bin/bash -c "/usr/bin/turnserver"
process_name = turn
